# backend/knapsack-backend.service
[Unit]
Description=Knapsack Tool Backend API
Documentation=https://github.com/your-repo
After=network.target mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/knapsack-tool/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=knapsack-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
